# src/main/resources/modbusrtu_exporter.service
[Unit]
Description=ModbusRTU Exporter

[Service]
Restart=always
User=pi
ExecStart=/usr/bin/java -Djava.util.logging.config.file=<path to logging.properties> -jar <path to modbusrtu_exporter.jar> <path to config.properties>

[Install]
WantedBy=multi-user.target